# cron.d/ddns-updater
*/30 * * * * root { test -x /usr/bin/ddns-updater || test -x /usr/local/bin/ddns-updater ; } && ddns-updater --syslog
